# templates/etc/cron.daily/apticron
0 14 * * * root /usr/sbin/apticron --cron
