# templates/etc/profile.d/alias.sh

alias sudo='sudo '
alias apt-get='echo "nope, please dont install stuff manually"'
