# templates/etc/network/if-up.d/add_gateway
#!/bin/bash

/sbin/ip route add default via {{ host_conf[env_name].public_network_gw }}
